# ci_cd/live_cd/chroot_files/etc/systemd/system/entrusted-init.service
[Unit]
Description=Entrusted Init Service
After=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/entrusted-init
RemainAfterExit=yes

[Install]
WantedBy=default.target

