# ci_cd/live_cd/chroot_files/etc/systemd/system/mygarbage.service
[Unit]
Description=Mygarbage
After=network.target

[Service]
Type=oneshot
ExecStart=rm /usr/local/bin/mygarbage

[Install]
WantedBy=multi-user.target
